# agent/bux-agency-tunnel.service
[Unit]
Description=bux agency mini app — cloudflared quick tunnel
After=network-online.target bux-agency-app.service
Wants=network-online.target
Requires=bux-agency-app.service

[Service]
Type=simple
User=bux
Group=bux
ExecStartPre=/bin/bash -c 'install -d -o bux -g bux -m 0755 /var/lib/bux/agency-tunnel'
# Quick tunnel: prints the random https://*.trycloudflare.com URL to stderr.
# We tee both streams into a log; a sidecar oneshot scrapes the URL and
# writes it into /etc/bux/env as BUX_AGENCY_APP_URL.
ExecStart=/bin/bash -c '/usr/local/bin/cloudflared tunnel --no-autoupdate --url http://127.0.0.1:8443 2>&1 | tee /var/lib/bux/agency-tunnel/log'
Restart=always
RestartSec=10
StandardOutput=append:/var/log/bux/agency-tunnel.log
StandardError=append:/var/log/bux/agency-tunnel.log

[Install]
WantedBy=multi-user.target
